# 20180401_script_check_ssh_connection/crontab
# executed each 3 hours
0 */3 * * * verify_ssh_connection_available.sh 
